# bot.conf
[program:kira_bot]
directory=/home/kira_bot
command=python3.8 main.py

autstart=true
autorestart=true

environment=HOME="/home/kira_bot",USER="root"

stderr_logfile=/home/kira_bot/logfile.log